# geradoresColect.service
[Unit]
Description=dataCenterColect
[Service]
Type=simple
ExecStart=/etc/init.d/geradoresColect start

[Install]
WantedBy=multi-user.target
